# resources/etc/profile.d/interactive_variables.sh
#export EDITOR=emacs
# export EDITOR="/usr/local/bin/code -w"
export EDITOR=~/bin/code-wait
# export LESSEDIT='mate -l %lm %f'
export LESSEDIT="/usr/local/bin/code -w --goto '%f\:%lm'"
